# configs/systemd/user/hyprpaper.service
[Unit]
Description=Hyprpaper wallpaper daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/hyprpaper
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
